# hosts/817caaa5-9ec3-4f70-9e52-a4f9baba849f/config/cron.d/infra-audit
*/5 * * * * root . /etc/environment && /opt/infra/scripts/audit_runner.sh >> /var/log/infra-audit.log 2>&1
