# hosts/002d015e-cd70-44e7-87d8-909c60ace2b3/config/cron.d/bdsmlr-metrics-rollup.j2
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

5 * * * * root /bin/bash /opt/services/config/metrics/rollup_metrics_hourly.sh >> /var/log/metrics-rollup.log 2>&1
